# deploy/splash-zero/debian/postinst
#!/bin/sh

set -e

systemctl daemon-reload >/dev/null 2>&1 || true
systemctl enable splash-serial.service >/dev/null 2>&1 || true

exit 0
